# claude-task-watcher.service
[Unit]
Description=监听 Obsidian Tasks 目录变化,文件一保存就立刻触发 Claude Code 任务扫描

[Service]
Type=simple
Environment=OBSIDIAN_VAULT=__OBSIDIAN_VAULT_PATH__
Environment=PATH=__RUNNER_PATH__
ExecStart=%h/.claude/skills/obsidian-task-runner/scripts/task-watcher.sh
Restart=on-failure
RestartSec=10

[Install]
WantedBy=default.target
