# claude-task-runner.timer
[Unit]
Description=定期兜底扫描一次 Obsidian 待办任务(事件触发已由 claude-task-watcher.service 负责,这里防止漏事件)

[Timer]
OnBootSec=5min
OnUnitActiveSec=30min
Persistent=true

[Install]
WantedBy=timers.target
