# claude-task-runner.service
[Unit]
Description=扫描 Obsidian Vault 并处理可执行的 Claude Code 任务(兜底轮询,由 timer 触发)

[Service]
Type=oneshot
Environment=OBSIDIAN_VAULT=__OBSIDIAN_VAULT_PATH__
# systemd --user 的 PATH 是干净的,不会读取 ~/.bashrc。
# install.sh 会尝试自动探测 claude/go 的路径填进来;如果探测不准,按 `which claude` 等结果手动改这一行。
Environment=PATH=__RUNNER_PATH__
ExecStart=%h/.claude/skills/obsidian-task-runner/scripts/task-runner-daemon.sh
