# omp-task-watcher.service
# This is a template. Run `otg install --systemd` to generate functional units.

[Unit]
Description=Obsidian Task Watcher — 监听文件变化,触发 OMP 处理

[Service]
Type=simple
Environment=OBSIDIAN_VAULT=__OBSIDIAN_VAULT_PATH__
Environment=PATH=__RUNNER_PATH__
ExecStart=%h/.local/bin/otg daemon
Restart=on-failure
RestartSec=10

[Install]
WantedBy=default.target
